# deploy/lodestar-coordinator@.service
# systemd template unit — one coordinator per tenant (instance-per-tenant).
#   sudo cp deploy/lodestar-coordinator@.service /etc/systemd/system/
#   # per-tenant env (unique FRAM_PORT + that tenant's log):
#   sudo install -d /etc/lodestar
#   printf 'FRAM_PORT=7801\nFRAM_LOG=/var/lib/lodestar/tenants/acme/claims.log\n' \
#     | sudo tee /etc/lodestar/acme.env
#   sudo systemctl enable --now lodestar-coordinator@acme
[Unit]
Description=Lodestar coordinator (tenant %i)
After=network.target

[Service]
Type=simple
User=lodestar
Group=lodestar
Environment=FRAM_HOME=/opt/fram
EnvironmentFile=/etc/lodestar/%i.env
ExecStart=/opt/fram/bin/fram-daemon ${FRAM_PORT} ${FRAM_LOG}
Restart=on-failure
RestartSec=2

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/lodestar

[Install]
WantedBy=multi-user.target
